# netlabelctl/90-netlabel.rules
SUBSYSTEM!="net", GOTO="netlabel_end"
RUN+="netlabel-config udev-action:$env{ACTION} device:$name"
LABEL="netlabel_end"
