# netlabelctl/netlabel.service.in
[Unit]
Description=NetLabel Kernel Configuration
After=systemd-modules-load.service
Before=network.target
Before=libvirtd.service
Before=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@prefix@/sbin/netlabel-config load
ExecStop=@prefix@/sbin/netlabel-config reset

[Install]
WantedBy=multi-user.target
